# deploy/systemd/anvil-agent.service
[Unit]
Description=Anvil Agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/anvil-agent
Environment=ANVIL_AGENT_HOST=127.0.0.1
Environment=ANVIL_AGENT_PORT=9090
Environment=INCUS_SOCKET=/var/lib/incus/unix.socket
ExecStart=/usr/local/bin/anvil-agent
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target
